# docker/crontab
*  *  *  *  * cd /home/geocertifica-api && php artisan schedule:run >> /dev/null 2>&1
